# crontab
# VMware IaaS定时任务配置
# 格式: 分 时 日 月 周 命令

# 每天凌晨2点检查过期虚拟机
0 2 * * * cd /app && python scheduler.py --check-expired >> /app/scheduler_logs/cron.log 2>&1

# 每天凌晨3点运行计费任务
0 3 * * * cd /app && python scheduler.py --daily-billing >> /app/scheduler_logs/cron.log 2>&1

# 每5分钟同步虚拟机状态
*/5 * * * * cd /app && python scheduler.py --sync-status >> /app/scheduler_logs/cron.log 2>&1

# 每小时清理过期会话
0 * * * * cd /app && python scheduler.py --cleanup-sessions >> /app/scheduler_logs/cron.log 2>&1

# 每天凌晨4点备份数据库
0 4 * * * cd /app && python scheduler.py --backup-db >> /app/scheduler_logs/cron.log 2>&1

# 确保cron日志目录存在
@reboot mkdir -p /app/scheduler_logs

# 空行确保crontab正确结束
